# nomad-autoscaler-setup-3/ansible/roles/prometheus/templates/prometheus.service.j2
[Unit]
Description=Prometheus
Documentation=https://prometheus.io/docs/
Wants=network-online.target
After=network-online.target

[Service]
User=prometheus
Group=prometheus
Type=simple
ExecStart=/usr/local/bin/prometheus --config.file=/etc/prometheus/prometheus.yml --storage.tsdb.path=/var/lib/prometheus
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
